# vm-systemd/chronyd.service.d/30_qubes.conf
[Unit]
ConditionPathExists=|/run/qubes-service/clocksync
ConditionPathExists=|/run/qubes-service/ntpd
After=qubes-sysinit.service
