# vm-systemd/getty@tty.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/run/qubes-service/getty@tty
After=qubes-sysinit.service
