# vm-systemd/user/pulseaudio.socket.d/30_qubes.conf
[Unit]
#ConditionNull=false
ConditionPathExists=/run/qubes-service/pulseaudio-vanilla-broken

