# vm-systemd/tor@default.service.d/30_qubes.conf
[Unit]
ConditionPathExists=!/run/qubes/this-is-templatevm
After=qubes-sysinit.service
