# vm-systemd/tinyproxy.service.d/30_not_needed_in_qubes_by_default.conf
[Unit]
ConditionPathExists=/run/qubes-service/tinyproxy
