# vm-systemd/netfilter-persistent.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/run/qubes-service/netfilter-persistent
After=qubes-sysinit.service
