# vm-systemd/network-manager.service.d/30_qubes.conf
# Disable sysinit version of network-manager (wheezy)
[Unit]
ConditionPathExists=!/run/qubes-service
After=qubes-sysinit.service
