# vm-systemd/NetworkManager-wait-online.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/run/qubes-service/network-manager
