# vm-systemd/anacron.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/run/qubes-service/crond

